# crontab
0 6 * * * /usr/bin/python3 /home/ubuntu/weather_script.py
